# env/master.conf
[program:master]
command=python voe/queue.py
directory=/home/ubuntu
autostart=true
autorestart=true
stdout_logfile=/home/ubuntu/master.log
redirect_stderr=true
stopsignal=QUIT
