# env/sink.conf
[program:master]
command=python sink.py
directory=/home/ubuntu
autostart=true
autorestart=true
stdout_logfile=/home/ubuntu/sink.log
redirect_stderr=true
stopsignal=QUIT
