# env/worker.conf
[program:worker]
command=python worker.py
directory=/home/ubuntu
autostart=true
autorestart=true
stdout_logfile=/home/ubuntu/worker.log
redirect_stderr=true
stopsignal=QUIT
